# etc/systemd/system/fix-disk-space.service
[Unit]
Description=Fix false "0 bytes free" after suspend/resume
After=suspend.target hibernate.target hybrid-sleep.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/fix-disk-space-after-resume.sh
RemainAfterExit=yes

[Install]
WantedBy=suspend.target hibernate.target hybrid-sleep.target
